# systemd/arch-auto-update.timer
[Unit]
Description=Run Lazy Arch Updater daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
